# pacotes/android-cmdtools/etc/profile.d/android-cmdtools.sh
export PATH="/opt/Android/Sdk/cmdline-tools/latest/bin:$PATH"
